# deploy/systemd/etazhi.service
[Unit]
Description=Etazhi lifts Flask app (Gunicorn)
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/etazhi-lifts
EnvironmentFile=/var/www/etazhi-lifts/.env
ExecStart=/var/www/etazhi-lifts/.venv/bin/gunicorn --workers 3 --bind 127.0.0.1:8000 wsgi:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

